# labs/linux_directories/files/Practice/Test/run/systemd/generator/local-fs.target.requires/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
What=/dev/disk/by-uuid/50a9826b-3a50-44d0-ad12-28f2056e9927
Where=/
Type=xfs
